# mirrror-mechanism/planet/etc.crontab
25    3 * * * /home/aptadm/bat/mirror-osm.planet

